# crontab
MAILTO=""
*/15 * * * * python ~/current/pypi_mirrors.py > ~/output/index.html
