# crontab
MAILTO=""
*/6 * * * * /home/dotcloud/env/bin/python /home/dotcloud/current/pypi_mirrors.py > /home/dotcloud/output/index.html
